# 8pit/mdocml/mdocml.cron
#!/bin/sh

local path="/usr/share/man"
local file=

for file in /usr/local/share/man/*; do
	if [ -e "${file}" -a "${file##*/}" != "mandoc.db" ]; then
		path="${path} /usr/local/share/man"
		break
	fi
done

echo $path
#exec /usr/sbin/makewhatis -a -T utf8 "${path}"
